# supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=/usr/sbin/php-fpm7.0 -F
autostart=true
autorestart=true
priority=999
username=www-data

[program:nginx]
command=/usr/sbin/nginx
username=www-data
autostart=true
autorestart=false
stdout_events_enabled=true
stderr_events_enabled=true

[program:php]
command=php artisan key:generate
username=www-data
autostart=true
autorestart=false
stdout_events_enabled=true
stderr_events_enabled=true
